# supervisord.conf
[supervisord]
nodaemon=true

[program:alfresco]
command=/alfresco/docker-entrypoint.sh